# stage2/01-sys-tweaks/files/rc.local
#!/bin/bash
source /boot/rpi_init
exit 0
